# actividad3/actividad3.service
[Unit]
Description=Saludo y fecha infinitos

[Service]
Type=simple
ExecStart=/usr/local/bin/saludo_script.sh

[Install]
WantedBy=multi-user.target